# discord-clan-bot.service
[Unit]
Description=Discord Clan Manager Bot
After=network.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/DiscordClanManagers
ExecStart=/usr/bin/python3 /home/botuser/DiscordClanManagers/main.py
Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true

StandardOutput=append:/home/botuser/DiscordClanManagers/bot.log
StandardError=append:/home/botuser/DiscordClanManagers/bot.log

[Install]
WantedBy=multi-user.target